# systemd/CcspLmLite.service
[Unit]
Description=CcspLMLite
After=CcspPandM.service CcspWiFiAgent.service
Requires=network-setup.service CcspWiFiAgent.service

[Service]
Type=forking
ExecStart=/usr/bin/CcspLMLite &
Restart=always

[Install]
WantedBy=multi-user.target

