# systemd/syseventd.service
[Unit]
Description=Sysevent-daemon
After=network-online.target
OnFailure=systemd-LastRebootReason.service

[Service]
ExecStart=/bin/sh -c '/lib/rdk/syseventd.sh;'
Type=forking

[Install]
WantedBy=multi-user.target



