# systemd/systemd-LastRebootReason.service
[Unit]
Description=LastRebootReason

[Service]
Type=forking
ExecStart=/bin/sh -c '/lib/rdk/LastRebootReason.sh;'
